# templates/bttrack/bttrack.service.erb
[Unit]
Description=Bittorent LDN Tracker
After=network.target auditd.service

[Service]
ExecStart=/usr/bin/bttrack --dfile <%= @dfile %> --logfile <%= @logfile %> --port <%= @port %> --nat_check <%= @nat_check %> --allowed_dir <%= @allowed_dir %> --ipv6_enabled <%= @ipv6_enabled %> --ipv6_binds_v4 <%= @ipv6_binds_v4 %> --bind <%= @bind %>
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=tracker.service
